# .config/systemd/user/chrome-tab-watchdog.timer
[Unit]
Description=Check Chrome tab memory every 30s

[Timer]
OnBootSec=60
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
