# .config/systemd/user/pnpm-update.service
[Unit]
Description=Update pnpm global packages
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
Environment=PATH=/home/adharsh/.local/share/fnm/aliases/default/bin:/usr/bin:/bin
ExecStart=/bin/bash -c 'pnpm update -g --latest'
# Rebuild @ladybugdb/core from source if prebuilt binary doesn't load (GLIBC mismatch)
ExecStartPost=/bin/bash -c '\
  PNPM_GLOBAL_DIR="$(dirname "$(pnpm root -g)")"; \
  LBUG_DIR=$(find "$PNPM_GLOBAL_DIR/.pnpm" -path "*/@ladybugdb/core/package.json" -not -path "*/lbug-source/*" -printf "%%h\n" 2>/dev/null | head -1); \
  if [ -n "$LBUG_DIR" ] && ! node -e "require(\"$LBUG_DIR/lbug_native.js\")" 2>/dev/null; then \
    (cd "$LBUG_DIR/lbug-source" && \
      CC=gcc-13 CXX=g++-13 cmake -B build/release \
        -DCMAKE_BUILD_TYPE=Release -DBUILD_NODEJS=TRUE \
        -DCMAKE_C_COMPILER=gcc-13 -DCMAKE_CXX_COMPILER=g++-13 . && \
      cmake --build build/release --config Release -j"$(nproc)") && \
    cp "$LBUG_DIR/lbug-source/tools/nodejs_api/build/lbugjs.node" "$LBUG_DIR/"; \
  fi'
