# .config/systemd/user/pnpm-store-prune.timer
[Unit]
Description=Weekly pnpm store prune

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
