# .config/systemd/user/pnpm-update.timer
[Unit]
Description=Daily pnpm global update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
